# unverified/rootfs/lib/systemd/system/custom_onboot.service
[Unit]
Description=Run custom onboot script
After=cpufrequtils.service


[Service]
Type=simple
ExecStart=/root/bin/onboot.bash

[Install]
WantedBy=graphical.target

