# unverified/rootfs/lib/systemd/system/user-zvC9-iptables.service
[Unit]
Description=user-zvC9 iptables

[Service]
Type=forking
ExecStart=/root/bin/user-zvC9-iptables.bash


[Install]
WantedBy=multi-user.target

